# init/batteryd.conf
description "Battery Daemon"
author "Vaibhav Verma <vaibhav2614@gmail.com>"

start on runlevel [2345]
stop on runlevel [016]

exec /usr/sbin/batteryd
